# telegram-bot.service
[Unit]
Description=Telegram Bot - Qarz Eslatmalari
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/alisher
Environment="PATH=/var/www/alisher/venv/bin"
ExecStart=/var/www/alisher/venv/bin/python3 /var/www/alisher/run_telegram_bot.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-bot

[Install]
WantedBy=multi-user.target
